# whisper-dictation.service
[Unit]
Description=Whisper Dictation (push-to-talk)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
ExecStart=%h/.local/bin/dictation --key f8
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
